# systemd/omen-backlight-sync.service
[Unit]
Description=HP OMEN Max real backlight sync workaround
Documentation=https://github.com/kcamporacosta/hp-omen-max-linux-brightness-fix
After=multi-user.target

[Service]
Type=simple
EnvironmentFile=/etc/omen-backlight/env
ExecStart=/usr/local/bin/omen-backlight-sync
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
